# deploy/guardian-monitor.timer
[Unit]
Description=Run Guardian API health monitor every 60 seconds

[Timer]
OnBootSec=60
OnUnitActiveSec=60
Unit=guardian-monitor.service

[Install]
WantedBy=timers.target
